# deploy/systemd/ec-consumer.service
[Unit]
Description=Extraction-Classification Kafka consumer
After=network.target kafka.service clickhouse-server.service
Requires=kafka.service
[Service]
User=tu
WorkingDirectory=/home/tu/sniff
ExecStart=/home/tu/sniff/.venv/bin/python -m integration.ec_consumer
Restart=always
RestartSec=5
[Install]
WantedBy=multi-user.target